# etc/update-motd.d/99-mining
#!/bin/bash
source /etc/mining.conf
source $TOOLS_PATH/FUNCTIONS
settings_load
echo "mining $DEFAULT_COIN"
echo ""
