# etc/update-motd.d/99-mining
#!/bin/bash
source /etc/mining.conf
$TOOLS_PATH/mu
echo ""
